# workspace/crontab/laradock
* * * * * laradock /bin/bash -c 'pgrep -xf "php /var/www/artisan queue:work --daemon" || php /var/www/artisan queue:work --daemon' >> /var/www/app/storage/logs/cron.log 2>&1
